# linux-install/customer-files/crond.conf
[program:crond]
autostart=true
autorestart=true
stopasgroup=true
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/cron.log
command=sudo /usr/sbin/crond -nsP
